# configs/db-int-01/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/sshd.out.log
stderr_logfile=/var/log/supervisor/sshd.err.log

[program:rsyslog]
command=/usr/sbin/rsyslogd -n -iNONE
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/rsyslog.out.log
stderr_logfile=/var/log/supervisor/rsyslog.err.log

[program:postgres]
command=/bin/bash -lc "su - postgres -s /bin/bash -c '/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf'"
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/postgres.out.log
stderr_logfile=/var/log/supervisor/postgres.err.log
